# .crontab
# brew nightly updates
0 3 * * *  brew update
0 4 * * *  brew upgrade
# npm
0 5 * * *  npm update npm -g
